# memo/autostart/fcitx.desktop
[Desktop Entry]
Exec=/usr/sbin/fcitx5
Icon=application-x-shellscript
Name=fcitx
Type=Application
Version=1.0
