# init/auto_setups/debian_script/rc.local
#!/bin/sh
/usr/local/bin/google-ime-skk
